# packaging/linux/deb/postinst
#!/bin/sh
set -e

SERVICE_NAME="pole-node.service"
CONFIG_DIR="/etc/pole"
DATA_DIR="/var/lib/pole"
LOG_DIR="/var/log/pole"

mkdir -p "$CONFIG_DIR" "$DATA_DIR" "$LOG_DIR"

systemctl daemon-reload || true
systemctl enable "$SERVICE_NAME" || true
systemctl restart "$SERVICE_NAME" || true

exit 0
